# misc/gutbackup@.timer
[Unit]
Description=GutBackup %i
After=network.target

[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
