# misc/gutbackup.service
[Unit]
Description=GutBackup Server mode
After=network.target

[Service]
ExecStart=/usr/bin/gutbackup server
