# misc/gutbackup@.service
[Unit]
Description=GutBackup %i
After=network.target

[Service]
ExecStart=/usr/bin/gutbackup backup %i
