# misc/gutbackup.timer
[Unit]
Description=GutBackup Server mode
After=network.target

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
